# templates/two_photo_api.supervisor.conf
[program:{{ app_name }}]
command={{ virtualenv_dir }}/bin/uwsgi --ini {{ project_dir }}/two_photo_api/uwsgi.ini
autostart=true
autorestart=true
stopwaitsecs=10
user=www-data
